# rootfs/etc/services.d/parsecd/run
#!/usr/bin/with-contenv /bin/bash

gosu parsec bash -c '/usr/bin/parsecd'

